# overlay.d/35oci-migration/usr/lib/systemd/system/zincati.service.d/010-oci-migration.conf
[Service]
StateDirectory=coreos-oci-migration
ExecStartPre=-/usr/libexec/coreos-oci-rebase
# wait for 3 min
TimeoutStartSec=300
