# overlay.d/35container-signing-migration/usr/lib/systemd/system/zincati.service.d/010-container-signing-migration.conf
[Service]
StateDirectory=coreos-container-signing-migration
ExecStartPre=-/usr/libexec/coreos-container-signing-migration
