# overlay.d/17fcos-container-signing/usr/lib/systemd/system/rpm-ostreed.service.d/coreos-containers-policy.conf
# Bind mount in our containers policy.json that configures
# how to do signature verification for quay.io/fedora/fedora-coreos
# Remove this override with:
#   mkdir /etc/systemd/system/rpm-ostreed.service.d
#   ln -s /dev/null /etc/systemd/system/rpm-ostreed.service.d/coreos-containers-policy.conf
[Service]
BindReadOnlyPaths=/usr/lib/coreos/coreos-containers-policy.json:/etc/containers/policy.json
